# misc/ops/systemd/role-janitor.timer
[Unit]
Description=Run Role Janitor hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
